# imagenes/subservicio7_logs/supervisord.conf
[supervisord]
nodaemon=true

[program:mongo]
command=/usr/bin/mongod

